# docker/all-in-one/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv sh
# Backend may still be loading BGE models on first start; nginx will retry
# upstream on its own, so we don't block here.
exec nginx -g 'daemon off;'
